# packaging/com.cat-break.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cat-break.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>__BINARY__</string>
    </array>

    <!-- Tweak the break threshold etc. here, then re-run install.sh. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>BREAK_MINUTES</key>
        <string>25</string>
        <key>CAT_SHOW_SECONDS</key>
        <string>8</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__LOG__</string>
    <key>StandardErrorPath</key>
    <string>__LOG__</string>
</dict>
</plist>
